# root/etc/services.d/nifi/run
#!/usr/bin/with-contenv sh

su ${USER} -c "${APP_HOME}/bin/nifi.sh run"